# deploy/systemd/uapk-gateway.service
[Unit]
Description=UAPK Gateway API Server
After=network.target uapk-chain.service
Wants=uapk-chain.service

[Service]
Type=simple
User=uapk
WorkingDirectory=/opt/uapk
EnvironmentFile=/opt/uapk/.env

# Run the gateway API (serves compiled instances)
ExecStart=/usr/bin/python3 -m uapk.cli run manifests/opspilotos.uapk.jsonld

Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=append:/var/log/uapk/gateway.log
StandardError=append:/var/log/uapk/gateway.log

[Install]
WantedBy=multi-user.target
